# somneo.service
[Unit]
Description=Submit temperature measurements to prometheus
After=network.target

[Service]
PrivateTmp=true
Type=simple
ExecStart=/usr/bin/somneo

[Install]
WantedBy=multi-user.target
